# bus/services/org.freedesktop.IBus.session.generic.service.in
[Unit]
Description=IBus Daemon for generic sessions
CollectMode=inactive-or-failed

# Never run in GNOME
Conflicts=gnome-session-initialized.target

[Service]
Type=dbus
# Only pull --xim in X11 session, it is done via Xwayland-session.d on Wayland
ExecStart=sh -c '@bindir@/ibus-daemon $IBUS_DAEMON_ARGS'
Restart=on-abnormal
BusName=org.freedesktop.IBus
TimeoutStopSec=5
Slice=session.slice
